# unicorn/upstart.conf.erb
description "API Unicorn application"

env USER=<%= user %>
env RAILS_ENV=production

start on runlevel [2]
stop on runlevel [016]

console owner

chdir <%= root_directory %>
exec su $USER -c "/usr/local/bin/bundle exec /usr/local/bin/unicorn -c ./config/unicorn/production.rb"

respawn

